# misc/scripts/worldrat-backup-postgresql.crontab
# m h  dom mon dow   command
0 0,18 * * * /usr/bin/worldrat-backup-postgresql
